# latest-email-server.service
[Unit]
Description=Latest Email Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=your-user
Group=your-group
WorkingDirectory=/path/to/latest-email-server
ExecStartPre=/path/to/go build -o app .
ExecStart=/path/to/latest-email-server/latest-email-server
Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=latest-email-server

[Install]
WantedBy=multi-user.target